# deploy/worker-host/ff-sims-discovery.service
[Unit]
Description=ff-sims discovery cron job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{SERVICE_USER}}
WorkingDirectory={{REPO_DIR}}/backend
EnvironmentFile=/etc/ff-sims-worker.env
ExecStart={{REPO_DIR}}/backend/cron -job=discovery -max-duration=50m
