# deploy/worker-host/ff-sims-lifetime-counts.service
[Unit]
Description=ff-sims lifetime-counts cron job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{SERVICE_USER}}
WorkingDirectory={{REPO_DIR}}/backend
EnvironmentFile=/etc/ff-sims-worker.env
# Type=oneshot services default to DefaultTimeoutStartSec (commonly 90s)
# before systemd kills them for "taking too long to start". A handful of
# COUNT queries finish in well under a minute, but this override matches
# ff-sims-discovery's for consistency and headroom.
TimeoutStartSec=10min
ExecStart={{REPO_DIR}}/backend/cron -job=lifetime-counts -max-duration=5m
