# deploy/worker-host/ff-sims-deploy.service
[Unit]
Description=ff-sims worker host self-update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=/root
WorkingDirectory={{REPO_DIR}}
ExecStart={{REPO_DIR}}/deploy/worker-host/deploy.sh
