# deploy/worker-host/ff-sims-lifetime-counts.timer
[Unit]
Description=Run ff-sims-lifetime-counts hourly

[Timer]
# Persistent=true is intentionally omitted: setup.sh's disable_sleep step
# already keeps the host from sleeping through a missed tick, so there's
# nothing for a catch-up run to make up for.
OnBootSec=3min
OnCalendar=hourly
Unit=ff-sims-lifetime-counts.service

[Install]
WantedBy=timers.target
